# geth-randomx.service
[Unit]
Description=Geth RandomX Node
Documentation=https://github.com/Aqui-oi/go-Ducros
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/go-Ducros
ExecStart=/home/ubuntu/go-Ducros/build/bin/geth \
  --datadir /home/ubuntu/go-Ducros/data-randomx \
  --networkid 33669 \
  --http \
  --http.addr "0.0.0.0" \
  --http.port 8545 \
  --http.api "eth,net,web3,miner,admin,debug" \
  --http.corsdomain "*" \
  --ws \
  --ws.addr "0.0.0.0" \
  --ws.port 8546 \
  --ws.api "eth,net,web3,miner,admin" \
  --mine \
  --miner.threads 4 \
  --miner.etherbase "0x0000000000000000000000000000000000000000" \
  --verbosity 3 \
  --maxpeers 50 \
  --log.file /home/ubuntu/go-Ducros/logs/geth.log

Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=geth-randomx

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
